# src/etc/cron.d/froxlor
PATH=/usr/local/sbin:/usr/local/bin:/sbin:/bin:/usr/sbin:/usr/bin
#
*/5 * * * * root /usr/bin/nice -n 5 /usr/bin/php -q /var/www/froxlor/bin/froxlor-cli froxlor:cron 'tasks' -q 1> /dev/null
0 0 * * * root /usr/bin/nice -n 5 /usr/bin/php -q /var/www/froxlor/bin/froxlor-cli froxlor:cron 'traffic' -q 1> /dev/null
5 0 * * * root /usr/bin/nice -n 5 /usr/bin/php -q /var/www/froxlor/bin/froxlor-cli froxlor:cron 'usage_report' -q 1> /dev/null
0 */6 * * * root /usr/bin/nice -n 5 /usr/bin/php -q /var/www/froxlor/bin/froxlor-cli froxlor:cron 'mailboxsize' -q 1> /dev/null
*/5 * * * * root /usr/bin/nice -n 5 /usr/bin/php -q /var/www/froxlor/bin/froxlor-cli froxlor:cron 'letsencrypt' -q 1> /dev/null
# Look for and purge old sessions every 30 minutes
09,39 * * * * root /usr/bin/nice -n 5 /usr/bin/php -q /var/www/froxlor/bin/froxlor-cli froxlor:php-sessionclean 1> /dev/null
